# TP-LINK/远程代码执行/CVE 2019-15060/TL-WR840N(TW)_V6.20_180726/squashfs-root/etc/inittab
::sysinit:/etc/init.d/rcS
ttyS1::respawn:/bin/sh
